# core/deploy/rc.local
#!/bin/bash
# Disable Pi LEDs for stealth
echo none > /sys/class/leds/ACT/trigger 2>/dev/null || true
echo 0    > /sys/class/leds/ACT/brightness 2>/dev/null || true
echo none > /sys/class/leds/PWR/trigger 2>/dev/null || true
echo 0    > /sys/class/leds/PWR/brightness 2>/dev/null || true
exit 0
